# supervisord.ini
[program:tornado_server]
directory = /mnt/e/PycharmProjects/flask-restful
command = python3 tornado_server.py
autostart = true
startsecs = 5
autorestart = true
startretries = 3
user = root
redirect_stderr = true
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 10
stopasgroup=true
killasgroup=true